# deploy/flyero-backup-failure@.service
[Unit]
Description=FLYERO Backup-Fehler im Systemjournal

[Service]
Type=oneshot
ExecStart=/usr/bin/logger -t flyero-backup "Backup fehlgeschlagen: %i. Pruefe journalctl -u flyero-backup.service."

