# deploy/flyero-notification-worker.timer
[Unit]
Description=FLYERO E-Mail-Benachrichtigungen zeitnah verarbeiten

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s
Persistent=true
Unit=flyero-notification-worker.service

[Install]
WantedBy=timers.target
